# VoiceFlow/Resources/app.getspit.relaunch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>app.getspit.relaunch</string>

    <!-- WATCHDOG (2026-07-10) — relança a app via `open` quando ela não está a
         correr. Substituiu o KeepAlive+BundleProgram anterior, que falhava por
         duas razões reais (diagnosticadas por dados):
           1. LWCR/EX_CONFIG: o launchd a spawnar o binário directamente valida
              um Launch Constraint que fica stale entre builds → spawn falhava.
           2. Conflito com o guard de instância única: o launchd criava uma 2ª
              instância que o guard matava com exit 0, e o launchd deixava de a
              gerir → morte por Jetsam nunca era recuperada.
         Este agent corre /bin/sh (sem constraint) e usa `open` (via LaunchServices,
         respeita a assinatura). `pgrep -x Spit` (nome exacto — NÃO -f, que faria
         self-match com esta própria linha de comando). O marcador .graceful_quit
         (escrito pela app em ⌘Q, removido no arranque) faz o watchdog respeitar
         uma saída deliberada e só relançar em morte anormal (crash/Jetsam). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>pgrep -x Spit &gt;/dev/null 2>&amp;1 || [ -f "$HOME/Library/Containers/app.getspit/Data/Library/Application Support/Spit/.graceful_quit" ] || open -g /Applications/Spit.app</string>
    </array>

    <!-- Corre no login e depois a cada 30s. Anti-loop natural: se a app arrancar
         e morrer logo, o intervalo limita a frequência. -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>30</integer>
</dict>
</plist>
